# ops/systemd/bcms-api-dev.service
[Unit]
Description=BCMS API stable local server
After=network-online.target rabbitmq-server.service snap.postgresql.postgresql.service bcms-opta-mount.service
Wants=network-online.target bcms-opta-mount.service
RequiresMountsFor=/mnt/opta-backups
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Desktop/bcms
EnvironmentFile=/home/ubuntu/Desktop/bcms/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=development
# LOW-INF-004 fix (2026-05-05): SKIP_AUTH dev için zorunlu; .env'den de
# yüklenebilir ama bu service "dev local-only" anlamı taşıyor.
# Production deploy bu file'ı kullanmaz; docker compose api servisi kullanır.
Environment=SKIP_AUTH=true
Environment=PORT=3000
ExecStartPre=/bin/sh -c '/usr/bin/fuser -k 3000/tcp || true'
ExecStart=/usr/bin/node /home/ubuntu/Desktop/bcms/apps/api/dist/server.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=default.target
