# ops/systemd/bcms-web-dev.system.service
[Unit]
Description=BCMS Web stable local server
After=network-online.target bcms-api-dev.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Desktop/bcms
EnvironmentFile=/home/ubuntu/Desktop/bcms/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=WEB_HOST=0.0.0.0
Environment=WEB_PORT=4200
Environment=API_PROXY_TARGET=http://127.0.0.1:3000
ExecStartPre=/bin/sh -c '/usr/bin/fuser -k 4200/tcp || true'
ExecStart=/usr/bin/node /home/ubuntu/Desktop/bcms/ops/scripts/bcms-web-static-server.mjs
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
