# ops/systemd/bcms-api-dev.system.service
[Unit]
Description=BCMS API stable local server
After=network-online.target rabbitmq-server.service snap.postgresql.postgresql.service bcms-opta-mount.service
Wants=network-online.target bcms-opta-mount.service
RequiresMountsFor=/mnt/opta-backups
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/Desktop/bcms
EnvironmentFile=/home/ubuntu/Desktop/bcms/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=development
Environment=SKIP_AUTH=true
Environment=PORT=3000
ExecStartPre=/bin/sh -c '/usr/bin/fuser -k 3000/tcp || true'
ExecStart=/usr/bin/node /home/ubuntu/Desktop/bcms/apps/api/dist/server.js
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
